# da-agent.service
[Unit]
Description=DA-Agent Uvicorn Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/DA-Agent
Environment="PATH=/root/.local/share/uv/python/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/root/DA-Agent"
ExecStart=/root/.local/bin/uv run uvicorn api.index:app --host 0.0.0.0 --port 8000 --timeout-keep-alive 300 --timeout-graceful-shutdown 300
Restart=always
RestartSec=10
StandardOutput=append:/root/DA-Agent/logs/da-agent.log
StandardError=append:/root/DA-Agent/logs/da-agent-error.log
SyslogIdentifier=da-agent

[Install]
WantedBy=multi-user.target