# da-agent.service
[Unit]
Description=DA-Agent Uvicorn Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/DA-Agent
Environment="PATH=/root/.local/share/uv/python/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/root/DA-Agent"
ExecStart=/root/.local/bin/uv run uvicorn api.index:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=da-agent

[Install]
WantedBy=multi-user.target